# Files/catalogue.service
[Unit]
Description = Catalogue Service

[Service]
User=roboshop
Environment=MONGO=true
Environment=MONGO_URL="mongodb://localhost:27017/catalogue"
WorkingDirectory=/app
ExecStart=/root/.nvm/versions/node/v21.2.0/bin/node server.js
Restart=always
SyslogIdentifier=catalogue

[Install]
WantedBy=multi-user.target